# deploy/systemd/personal-invest-backend.service
[Unit]
Description=Personal Invest Backend API
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/root/remote/personal-invest
Environment=ENV_FILE=.env.server
Environment=PYTHONUNBUFFERED=1
Environment=PYTHON_RUNTIME_BIN=/root/remote/personal-invest/.venv/bin/python
Environment=PATH=/root/.local/bin:/root/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/bin
ExecStart=/root/remote/personal-invest/scripts/backend_prod.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT

# Optional hardening. Enable after confirming paths and permissions.
# NoNewPrivileges=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target
